# deploy/systemd/innbucks-bot.service
[Unit]
Description=InnBoost Telegram Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=innbucks
WorkingDirectory=/opt/innbucks
EnvironmentFile=/opt/innbucks/.env
ExecStart=/opt/innbucks/bin/bot
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
